# system-monitor.service
[Unit]
Description=Custom System Monitor Tool
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/system-monitor
Restart=on-failure

[Install]
WantedBy=multi-user.target
